# zulip_tumcsbot.conf
[program:tumcsbot]
command=docker-compose -f /home/tumcsbot/tumcsbot/docker-compose.yml up
priority=999                   ; the relative start priority (default 999)
autostart=true                 ; start at supervisord start (default: true)
autorestart=false              ; whether/when to restart (default: unexpected)
stopsignal=TERM                ; signal used to kill process (default TERM)
stopwaitsecs=10                ; max num secs to wait b4 SIGKILL (default 10)
user=tumcsbot                  ; setuid to this UNIX account to run the program
redirect_stderr=true           ; redirect proc stderr to stdout (default false)
stdout_logfile=/home/tumcsbot/tumcsbot/tumcsbot.log
